# docker/images/web/data/etc/service/httpd/run
#!/usr/bin/env svrun
# vim: ai ts=2 sts=2 et sw=2 ft=ruby

Dir.chdir('/var/www/localhost') do
  service(['bundle',
           'exec',
           'rake',
           'serve',
           'serve_port=80',
           "serve_storage=/var/serve"],
          user: :'www-data',
          group: 'www-data').call
end

# Local Variables:
# mode: ruby
# End:
