# docker/images/assets/data/etc/service/httpd/run
#!/usr/bin/env svrun
# vim: ai ts=2 sts=2 et sw=2 ft=ruby

Dir.chdir('/var/www/localhost') do
  service(['darkhttpd',
           '/var/www/localhost',
           '--chroot',
           '--no-listing',
           '--port', '80',
           '--uid', ENV.fetch('VOLUME_UID'),
           '--gid', ENV.fetch('VOLUME_GID'),
           '--log', '/var/log/darkhttpd/access.log']).call
end

# Local Variables:
# mode: ruby
# End:
